# examples/com.example.screenshotcapture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Auto-capture Mac screenshots into Noggin's inbox.

  1. Copy screenshot-capture.sh somewhere permanent, edit its variables,
     and `chmod +x` it.
  2. Edit the two paths below (script location + your screenshot folder).
  3. Copy this file to ~/Library/LaunchAgents/ and run:
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.example.screenshotcapture.plist

  WatchPaths fires the script whenever the folder changes; the script itself
  filters down to actual macOS screenshots, so other files are left alone.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.screenshotcapture</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/full/path/to/screenshot-capture.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/yourname/Desktop</string>
    </array>
</dict>
</plist>
